# deployment/systemd/visa-bulletin.service
# Systemd service file for Visa Bulletin Dashboard
# Location on server: /etc/systemd/system/visa-bulletin.service
#
# Installation:
#   sudo cp visa-bulletin.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable visa-bulletin
#   sudo systemctl start visa-bulletin
#
# Management:
#   sudo systemctl status visa-bulletin
#   sudo systemctl restart visa-bulletin
#   sudo journalctl -u visa-bulletin -f

[Unit]
Description=Visa Bulletin Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/visa_bulletin
Environment="DJANGO_SETTINGS_MODULE=django_config.settings"

# Load secrets from .env file (secure)
# Create this file on server: /opt/visa_bulletin/.env
# Content: DJANGO_SECRET_KEY=...
EnvironmentFile=/opt/visa_bulletin/.env

# Analytics: GoatCounter privacy-focused tracking
# This can also be overridden in .env
Environment="ANALYTICS_SCRIPT=<script data-goatcounter=\"https://vyakunin.goatcounter.com/count\" async src=\"https://gc.zgo.at/count.js\"></script>"

# Gunicorn: 3 workers (2 * CPU + 1) for better concurrency
ExecStart=/opt/visa_bulletin/venv/bin/gunicorn --workers 3 --threads 2 --bind 0.0.0.0:8000 --timeout 120 --max-requests 1000 --max-requests-jitter 50 django_config.wsgi:application
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
